# root/usr/lib/systemd/system/ejabberd-create-dhparam.service
[Unit]
Description=Create the dhparam group for ejabberd

[Service]
Type=simple
ExecStart=/usr/libexec/nethserver/ejabberd-create-dhparam
